# systemd/labrad-node.service.d/labrad-node.conf
[Service]
Environment="LABRADHOST=localhost"
Environment="LABRADPASSWORD=698gang"
Environment="LABRADNODE=appa"
Environment="LABRAD_TLS=off"

